# supervisor.ini
; /etc/supervisord.d/filipeelore.ini

[program:filipeelore]
command=/root/.pyenv/versions/filipeelore.love/bin/gunicorn -k uvicorn.workers.UvicornWorker main:app --bind=0.0.0.0:5000 --workers=1
directory=/var/www/filipeelore.love
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/filipeelore.love.err.log
stdout_logfile=/var/log/filipeelore.love.out.log